# scripts/systemd/tuyaplug.service
[Unit]
Description=

[Service]
ExecStart=/usr/bin/bash start.sh
WorkingDirectory=/home/reed/proj/tuyaplug

[Install]
WantedBy=default.target
